# cmd/agent_system/package/authentik-agent.service
[Unit]
Description=authentik Agent
After=network.target

[Service]
Restart=always
ExecStart=/usr/bin/ak-agent

[Install]
WantedBy=multi-user.target
