# cmd/agent_system/package/postinst.sh
#!/usr/bin/env bash
systemctl enable --now authentik-session-manager
