# cmd/agent_system/package/linux/postinst.sh
#!/usr/bin/env bash
set -euo pipefail

case "$1" in
    configure)
    # creating _ak-agent group if he isn't already there
    if ! getent group _ak-agent >/dev/null; then
            addgroup --system --force-badname _ak-agent
    fi
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

systemctl --system daemon-reload >/dev/null || true
systemctl enable --now 'ak-sysd.service' >/dev/null || true

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
